# unicorn_upstart.conf
description "Unicorn Daemon for My Rails App"

start on runlevel [2]
stop on runlevel [016]

console owner

exec /usr/local/rvm/bin/r192_unicorn -c /var/www/test_app/config/unicorn.rb

respawn
